# attestation/pckid-service/pckid-retrieve-tool-service_1.0-1_amd64/DEBIAN/postinst
systemctl enable pckid_retrieve_tool_service.service
